# src/resources/cirrus.service
[Unit]
Description=Cirrus Backup

[Service]
Restart=always
ExecStart={cirrus_binary} daemon

[Install]
WantedBy=default.target
