# selfinstaller/examples/scheduledjob.timer
[Unit]
Description=Regularly print a supportive message

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
