# selfinstaller/examples/scheduledjob.service
[Unit]
Description=Sample job that prints a supportive message

[Service]
Type=oneshot
ExecStart={{executable}} run "Hey there, you are great!"
